# ltr-fix_0.0.1-1_amd64/DEBIAN/prerm
#!/bin/bash

deb-systemd-helper stop ltr-fix.service >/dev/null || true
deb-systemd-helper disable ltr-fix.service >/dev/null || true
systemctl --system daemon-reload >/dev/null || true
